# files/telegram-listener.conf
[supervisord]
loglevel=info

[inet_http_server]
port = *:9001

[program:cakesbot]
command=/var/opt/babyml/cakesbot.py
startsecs=15
directory=/var/opt/babyml