# files/telegram-listener.conf
[supervisord]
loglevel=info

[inet_http_server]
port = *:9001

[program:telegram_listener]
command=/var/opt/baby_predict/listener.py
startsecs=15
directory=/var/opt/baby_predict